# linuxService/goldmarket.service
[Unit]
Description=GoldMarket Comprehensive Intelligence System - Backend API
Documentation=https://github.com/zw0269/GoldMarketComprehensiveIntelligenceSystem
# 等待网络完全就绪（含 DNS 解析）后再启动，确保数据采集正常
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/GoldMarketComprehensiveIntelligenceSystem

# 加载 .env 中的 API Keys 和配置
EnvironmentFile=/var/www/GoldMarketComprehensiveIntelligenceSystem/.env
Environment=NODE_ENV=production
Environment=HOME=/home/ubuntu
Environment=PATH=/home/ubuntu/.nvm/versions/node/v24.14.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# 使用 tsx 直接运行 TypeScript，无需手动编译
# 若已执行 npm run build，可改为：node dist/index.js
ExecStart=/home/ubuntu/.nvm/versions/node/v24.14.1/bin/node \
    /var/www/GoldMarketComprehensiveIntelligenceSystem/node_modules/.bin/tsx \
    src/index.ts

# 崩溃后自动重启，最多 10 秒内重试
Restart=always
RestartSec=10
# 1 分钟内启动失败超过 5 次则停止尝试
StartLimitInterval=60
StartLimitBurst=5

# 资源限制（防止单进程占用过多）
LimitNOFILE=65536

# 日志输出到 systemd journal（用 journalctl -u goldmarket -f 查看）
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goldmarket

[Install]
WantedBy=multi-user.target
